# etc/systemd/system/gerrit-xxx.service
[Unit]
Description=xxx Gerrit
Requires=gerrit@xxx.service
[Service]
Type=oneshot
ExecStart=/bin/true
[Install]
WantedBy=multi-user.target
