# etc/systemd/system/git-sync.path
[Unit]
Description=git-sync scanner
Requires=git-sync-pre.service
[Path]
PathExistsGlob=/run/git-sync/*.trigger
[Install]
WantedBy=multi-user.target
