# etc/systemd/system/git-sync.service
[Unit]
Description=git-sync
[Service]
Type=oneshot
User=git-sync
ExecStart=/home/git-sync/git-sync/bin/git-sync.service
